# musiccards.service
[Unit]
Description=MusicCards
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/music-cards
ExecStart=/usr/bin/python -u /home/pi/music-cards/box.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target